# root/etc/s6-overlay/s6-rc.d/svc-pulseaudio/run
#!/usr/bin/env bash
set -euo pipefail

DISPLAY=${DISPLAY:-:0}
XSOCK="/tmp/.X11-unix/X${DISPLAY#:}"

echo "[init-sunshine] Waiting for X at $DISPLAY ..."
for i in {1..60}; do
  [ -S "$XSOCK" ] && break
  sleep 1
done
[ -S "$XSOCK" ] || { echo "[init-sunshine] X socket not found, exiting."; exit 1; }

# Ensure runtime and config dirs exist
mkdir -p /run/user/1000 /config/.config/pulse /config/.config/sunshine
chown -R abc:abc /run/user/1000 /config/.config

# ---- PULSEAUDIO INITIALIZATION ----
if ! pgrep -u abc pulseaudio >/dev/null 2>&1; then
  echo "[init-sunshine] Starting PulseAudio for abc..."
  s6-setuidgid abc pulseaudio --start --exit-idle-time=-1 || true
else
  echo "[init-sunshine] PulseAudio already running"
fi

# Export environment
export DISPLAY
export XDG_RUNTIME_DIR=/run/user/1000
export HOME=/config
export SUNSHINE_CONFIG_DIR="/config/.config/sunshine"

# ---- CLEANUP ----
cleanup() {
  echo "[init-sunshine] Stopping Sunshine..."
  pkill -u abc -f "/usr/bin/sunshine" || true
  pkill -u abc -f "pulseaudio" || true
}
trap cleanup SIGTERM SIGINT EXIT

# ---- START SUNSHINE ----
echo "[init-sunshine] Starting Sunshine as user abc"
exec s6-setuidgid abc /usr/bin/sunshine
