# root/etc/s6-overlay/s6-rc.d/init-sunshine/finish
#!/usr/bin/execlineb -P
s6-echo "sunshine: exited with code $1, restarting if needed..."
