# root/etc/s6-overlay/s6-rc.d/init-sunshine/run
#!/bin/bash
set -e

echo "sunshine: waiting for X at ${DISPLAY:-:0} ..."
for i in {1..60}; do
  if [ -S /tmp/.X11-unix/X0 ]; then
    echo "sunshine: X socket found, proceeding..."
    break
  fi
  sleep 5
done

if [ ! -S /tmp/.X11-unix/X0 ]; then
  echo "sunshine: X socket not found after timeout, exiting..."
  exit 1
fi

echo "sunshine: waiting for PulseAudio to initialize..."
sleep 2

echo "sunshine: PulseAudio initialized, starting Sunshine..."
exec /usr/bin/sunshine --config-dir /config/sunshine
