# root/etc/s6-overlay/s6-rc.d/init-sunshine/run
#!/usr/bin/with-contenv bash
set -euo pipefail

# Disable system/session dbus so no activation attempts occur
unset DBUS_SESSION_BUS_ADDRESS
unset DBUS_SYSTEM_BUS_ADDRESS

DISPLAY=${DISPLAY:-:0}
XSOCK="/tmp/.X11-unix/X${DISPLAY#:}"

echo "sunshine: waiting for X at $DISPLAY ..."
for i in {1..60}; do
  if [ -S "$XSOCK" ]; then
    echo "sunshine: X socket found after $((i*5))s"
    break
  fi
  sleep 5
done

if [ ! -S "$XSOCK" ]; then
  echo "sunshine: X socket not found after timeout, exiting"
  exit 1
fi

# Ensure runtime and config dirs
mkdir -p /run/user/1000 /config/.config/pulse /config/.config/sunshine
chown -R abc:abc /run/user/1000 /config/.config

# Start PulseAudio
if ! pgrep -u abc pulseaudio >/dev/null 2>&1; then
  echo "sunshine: starting PulseAudio"
  s6-setuidgid abc pulseaudio --start --exit-idle-time=-1 || true
else
  echo "sunshine: PulseAudio already running"
fi

# Export environment
export DISPLAY=:0
export XDG_RUNTIME_DIR=/run/user/1000
export HOME=/config

# Start Sunshine
echo "sunshine: launching Sunshine as user abc..."
exec s6-setuidgid abc /usr/bin/sunshine
