# root/etc/services.d/sunshine/finish
#!/usr/bin/execlineb -P
# If Sunshine exits, log and let s6 restart it (default behavior)
s6-echo "sunshine: exited with code $1, restarting if needed..."
